# email/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=10s
OnUnitActiveSec=10m
Unit=mbsync.service

[Install]
WantedBy=timers.target
