# backup/backup.service
[Unit]
Description=Backup

[Service]
Type=oneshot
ExecStart=/home/ccabessa/.local/bin/backup.sh
